# gatherClear.cron
#!/bin/sh

cd /root/ib/
./gather cron
./clear
